# playbooks/roles/infra/kopia_server/templates/kopia-verify.timer.j2
[Unit]
Description=Kopia repository verification (daily 5% sample)

[Timer]
OnCalendar=daily
RandomizedDelaySec=7200
Persistent=true

[Install]
WantedBy=timers.target
