# playbooks/roles/services/inpx_web/files/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:inpx-web]
command=sh -c 'INPX=$(find /downloads -maxdepth 1 -name "*.inpx" | head -1); exec node /srv/inpx-web/server --data-dir=/data/inpx-web --lib-dir=/downloads/Flibusta.Net ${INPX:+--inpx=$INPX}'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:liberama]
command=node /srv/liberama/server --app-dir=/data/liberama
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
